# pm-system/systemd/popcrm-apply-ignore-rules.timer
[Unit]
Description=Run POP CRM ignore-rule sweep every 6 hours

[Timer]
OnCalendar=*-*-* 00,06,12,18:35:00
Persistent=true
Unit=popcrm-apply-ignore-rules.service

[Install]
WantedBy=timers.target
