# pm-system/systemd/popcrm-summarize.service
[Unit]
Description=POP CRM opportunity summary refresh
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=ai
WorkingDirectory=/worksp/directus
EnvironmentFile=/home/ai/.crm-worker.env
ExecStart=/usr/bin/node /worksp/directus/pm-system/crm-worker.mjs summarize
