# pm-system/systemd/popcrm-contact-sync.timer
[Unit]
Description=Run POP CRM contact sync daily

[Timer]
OnCalendar=*-*-* 02:10:00
Persistent=true
Unit=popcrm-contact-sync.service

[Install]
WantedBy=timers.target
